# my-sci-app/my-daemon.conf
[program:cat]
command = /bin/cat
autostart = true
autorestart = true
