# base-sci-app/sshd.conf
[program:sshd]
command = /usr/sbin/sshd -D
stdout_logfile = /var/log/supervisor/sshd_out.log
stderr_logfile = /var/log/supervisor/sshd_err.log
autorestart = true
